# docker/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/docker-utilities/artisan schedule:run >> /dev/null 2>&1
